# src/ansible02/roles/install_consul_client/templates/envoy.service.j2
[Unit]
Description=Consul Envoy
After=network.target

[Service]
ExecStart=/usr/local/bin/envoy -c /etc/consul.d/envoy.yaml
Restart=always
User=consul

[Install]
WantedBy=multi-user.target