# src/ansible02/roles/install_hotels_service/templates/hotel-service.service.j2
[Unit]
Description=Hotel Service
After=network.target consul.service

[Service]
EnvironmentFile=/etc/hotel-service.env
ExecStart=/usr/bin/java -jar /opt/hotel-service.jar
User=hoteluser
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target